# cron.d/elasticsearch-backup
0 4 * * * root /usr/local/bin/elasticsearch-backup >> /var/log/ext/elasticsearch-backup.log 2>&1
