# cron.d/elasticsearch-backup
0 4 * * * root /usr/local/bin/elasticsearch-backup
